# contrib/systemd/signalforge-agent.service
[Unit]
Description=SignalForge agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SIGNALFORGE_AGENT_USER__
WorkingDirectory=__SIGNALFORGE_AGENT_WORKDIR__
EnvironmentFile=__SIGNALFORGE_AGENT_ENV_FILE__
ExecStart=__SIGNALFORGE_AGENT_BUN__ run src/cli.ts run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
